# deploy/reel.service
[Unit]
Description=MEDIA! Media Server
After=network.target

[Service]
Type=simple
User=reel
Group=reel
WorkingDirectory=/opt/media-app
ExecStart=/bin/bash scripts/start-prod.sh
Restart=on-failure
RestartSec=10
Environment=NODE_ENV=production
Environment=PATH=/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
